# adguard/rootfs/etc/s6-overlay/s6-rc.d/doorbell/run
! /bin/bash
touch /data/adguard/data/querylog.json
sed -i "s/size_memory: 1000/size_memory: 0/g" /data/adguard/AdGuardHome.yaml
#s6-svc -r /var/run/service/adguard
tail -fn0 /data/adguard/data/querylog.json | while read;do echo "$REPLY" | grep -q "ali-alarm-eu.oss-eu-central-1.aliyuncs.com" &&
curl -k --header "Content-Type: application/json" \
  --request POST \
  --data '{"pressed": true}' \
  http://homeassistant.local:8123/api/webhook/doorbell_button_1423! ; done
